# ansible/roles/dandydash-api/templates/dandydash-api.service.j2
[Unit]
Description=DandyDash API - {{ domain }}
After=network.target postgresql.service

[Service]
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment="PATH={{ app_dir }}/venv/bin"
EnvironmentFile={{ app_dir }}/.env
ExecStart={{ app_dir }}/venv/bin/gunicorn \
    --workers 4 \
    --bind 127.0.0.1:5000 \
    --access-logfile {{ log_dir }}/access.log \
    --error-logfile {{ log_dir }}/error.log \
    run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
